# deploy/movemarias.service
[Unit]
Description=Move Marias Django Application
After=network.target redis.service
Wants=redis.service

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=movemarias
WorkingDirectory=/var/www/movemarias
Environment=PATH=/var/www/movemarias/venv/bin
EnvironmentFile=/var/www/movemarias/.env

# Gunicorn configuration
ExecStart=/var/www/movemarias/venv/bin/gunicorn \
    --config /var/www/movemarias/deploy/gunicorn.conf.py \
    movemarias.wsgi:application

# Reload configuration
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/movemarias/media /var/log/movemarias /var/run/movemarias

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=movemarias

[Install]
WantedBy=multi-user.target
